# watch.conf
[program:watch]
command=python /home/jinank94/GitProjects/watchguard-manager/wgtool.py &
autostart=true
autorestart=true
stderr_logfile=/var/log/watch.err.log
stdout_logfile=/var/log/watch.out.log
